# wenet_tx.service
[Unit]
Description=wenet_tx
After=basic.target

[Service]
# Update this path if the wenet files are not in the pi user's home directory!
ExecStart=/home/pi/wenet/start_tx.sh
Restart=always
RestartSec=30
# This one too!
WorkingDirectory=/home/pi/wenet/

# NOTE - Wenet unfortunately needs to be run at root, else we can't easily set the system time.
#User=pi
SyslogIdentifier=wenet_tx

[Install]
WantedBy=multi-user.target